# extra/systemd-services/indexer-web-server-watcher.path
[Path]
PathModified=/home/indexer/indexer/

[Install]
WantedBy=multi-user.target
